# kiro2cc-proxy.service
[Unit]
Description=kiro2cc-proxy Anthropic API Reverse Proxy
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/kiro2cc-proxy

# 从 .env 文件加载配置（文件不存在时忽略）
EnvironmentFile=-/opt/kiro2cc-proxy/.env

ExecStart=/opt/kiro2cc-proxy/kiro2cc-proxy
Restart=always
RestartSec=10

# 日志写入 journald（用 journalctl -u kiro2cc-proxy -f 查看）
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
